# conf.d/emulator.conf
[program:emulator]
command=/opt/android-sdk/emulator/emulator -avd automotive_emulator -no-audio -verbose -show-kernel -skin 1080x600
autorestart=true
